# provision/service/apache2.conf
[program:apache2]
command=/usr/sbin/apachectl -e info -DFOREGROUND
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true

